# projects/splunk/roles/uf/templates/splunk_init
#!/bin/sh
#
# /etc/init.d/splunk-{{ splunkforwarder_role }}
# init script for Splunk.
# generated by 'Jay Slay....and ansible'
#
# chkconfig: 2345 90 60
# description: Splunk {{ splunkforwarder_role }} service
#
RETVAL=0

. /etc/init.d/functions

splunk_start() {
  echo Starting Splunk...
  "{{ splunkforwarder_home }}/bin/splunk" start --answer-yes --accept-license
  RETVAL=$?
  [ $RETVAL -eq 0 ] && touch /var/lock/subsys/splunk-{{ splunkforwarder_role }}
}
splunk_stop() {
  echo Stopping Splunk...
  "{{ splunkforwarder_home }}/bin/splunk" stop 
  RETVAL=$?
  [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/splunk-{{ splunkforwarder_role }}
}
splunk_restart() {
  echo Restarting Splunk...
  "{{ splunkforwarder_home }}/bin/splunk" restart 
  RETVAL=$?
  [ $RETVAL -eq 0 ] && touch /var/lock/subsys/splunk-{{ splunkforwarder_role }}
}
splunk_status() {
  echo Splunk status:
  "{{ splunkforwarder_home }}/bin/splunk" status 
  RETVAL=$?
}
case "$1" in
  start)
    splunk_start
    ;;
  stop)
    splunk_stop
    ;;
  restart)
    splunk_restart
    ;;
  status)
    splunk_status
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac

exit $RETVAL
